# deploy/ec2/saathi-api.service
[Unit]
Description=Saathi FastAPI (Gunicorn + Uvicorn workers)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/saathi/backend
EnvironmentFile=-/opt/saathi/backend/.env
# So ec2-user can `pip install -r requirements.txt` without sudo after someone used `sudo pip` or mixed owners.
ExecStartPre=/bin/chown -R ec2-user:ec2-user /opt/saathi/backend/venv
ExecStart=/opt/saathi/deploy/ec2/run-gunicorn.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
